# quantumfond-backend.service
[Unit]
Description=QuantumFond Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=quantumfond
Group=quantumfond
WorkingDirectory=/opt/quantumfond/backend
Environment="PATH=/opt/quantumfond/backend/venv/bin"
EnvironmentFile=/opt/quantumfond/backend/.env
ExecStart=/opt/quantumfond/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
