# quantum-safety-telemetry.service
[Unit]
Description=Quantum Trader - Safety Telemetry Exporter (P1)
Documentation=https://prometheus.io/docs/instrumenting/clientlibs/
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/safety-telemetry.env
Environment="PATH=/opt/quantum/venvs/safety-telemetry/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/quantum/venvs/safety-telemetry/bin/python3 /home/qt/quantum_trader/microservices/safety_telemetry/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-safety-telemetry

# Restart policy
Restart=always
RestartSec=3

# Resource limits (lightweight)
MemoryHigh=50M
MemoryMax=64M

# Security hardening (read-only to Redis)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
