# quantum-strategic-evolution.service
[Unit]
Description=Quantum Strategic Evolution Engine
Documentation=https://github.com/binyaminsemerci-ops/quantum_trader
After=network.target quantum-core.target
PartOf=quantum-ai.target
Requires=quantum-core.target

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
Environment="PATH=/opt/quantum/venvs/ai-engine/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/qt/quantum_trader"
ExecStart=/opt/quantum/venvs/ai-engine/bin/python -m microservices.strategic_evolution.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-strategic-evolution

[Install]
WantedBy=quantum-ai.target
