# ops/systemd/quantum-state.target
[Unit]
Description=State tracking (positions, portfolio)
Requires=quantum-core.target
After=quantum-core.target

[Install]
WantedBy=quantum-trader.target
