# ops/systemd/quantum-emergency-exit-worker.service
[Unit]
Description=Quantum Trader Emergency Exit Worker
Documentation=https://github.com/quantum-trader/services/emergency_exit_worker
After=network.target redis-server.service
# CRITICAL: This service must always be running when trading is active
StartLimitIntervalSec=0

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/qt/quantum_trader
EnvironmentFile=/home/qt/quantum_trader/.env

# Main execution
# NOTE: Uses ai-client-base venv (/home/qt/quantum_trader/venv intentionally absent)
ExecStart=/opt/quantum/venvs/ai-client-base/bin/python \
    /home/qt/quantum_trader/services/emergency_exit_worker/emergency_exit_worker.py

# NO automatic restart - failures require manual inspection
Restart=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-emergency-exit-worker

# Resource limits
MemoryMax=256M
CPUQuota=25%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
