# ops/systemd/quantum-risk_safety.service
[Unit]
Description=Risk management and safety checks
Requires=quantum-core.target quantum-state.target
After=quantum-core.target quantum-state.target

[Service]
Type=simple
User=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/risk_safety.env
ExecStart=/home/qt/.local/bin/uv run python microservices/risk_safety/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=quantum-brains.target
