# ops/systemd/quantum-core.target
[Unit]
Description=Core infrastructure (Redis, market data)



[Install]
WantedBy=quantum-trader.target
