# ops/systemd/quantum-exposure_balancer.service
[Unit]
Description=Exposure balancing across positions
Requires=quantum-core.target quantum-state.target quantum-brains.target quantum-ai.target
After=quantum-core.target quantum-state.target quantum-brains.target quantum-ai.target

[Service]
Type=simple
User=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/exposure_balancer.env
ExecStart=/home/qt/.local/bin/uv run python microservices/exposure_balancer/service.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=quantum-exec.target
