# ops/systemd/quantum-obs.target
[Unit]
Description=Observability and logging
Requires=quantum-core.target
After=quantum-core.target

[Install]
WantedBy=quantum-trader.target
