# ops/systemd/quantum-ai.target
[Unit]
Description=AI and RL components
Requires=quantum-state.target
After=quantum-state.target

[Install]
WantedBy=quantum-trader.target
