# ops/systemd/quantum-rl_training.service
[Unit]
Description=RL-based position sizing agent
Requires=quantum-core.target quantum-state.target
After=quantum-core.target quantum-state.target

[Service]
Type=simple
User=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/rl_training.env
ExecStart=/home/qt/.local/bin/uv run python microservices/rl_training/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=quantum-ai.target
