# deployment/systemd/quantum-governor.service
[Unit]
Description=Quantum Trading P3.2 Governor - Fund-Grade Limits + Auto-Disarm
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/governor.env
ExecStart=/usr/bin/python3 /home/qt/quantum_trader/microservices/governor/main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-governor

[Install]
WantedBy=multi-user.target
