# deploy/systemd/quantum-capital-efficiency.service
[Unit]
Description=Quantum P3.1 Capital Efficiency Brain
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=/etc/quantum/capital-efficiency.env
ExecStart=/usr/bin/python3 /home/qt/quantum_trader/microservices/capital_efficiency/main.py
Restart=always
RestartSec=2

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/qt/quantum_trader /tmp
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
