# ops/systemd/native/quantum-execution.service
[Unit]
Description=Quantum Trader Execution Service
Documentation=https://github.com/binyaminsemerci-ops/quantum_trader
After=quantum-core.target quantum-ai.target
PartOf=quantum-exec.target
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/qt/quantum_trader
Environment="PYTHONPATH=/home/qt/quantum_trader"
EnvironmentFile=/etc/quantum/execution.env
ExecStart=/opt/quantum/venvs/execution/bin/uvicorn microservices.execution.service_v2:app --host 0.0.0.0 --port 8003 --log-level info

# Restart policy
Restart=always
RestartSec=2

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Resource limits
MemoryMax=1G
MemoryHigh=900M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-execution

[Install]
WantedBy=quantum-exec.target
