# ops/systemd/native/quantum-ai-engine.service
[Unit]
Description=Quantum Trader AI Engine Service
Documentation=https://github.com/binyaminsemerci-ops/quantum_trader
After=quantum-core.target
PartOf=quantum-ai.target
Requires=redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/qt/quantum_trader
Environment="PYTHONPATH=/home/qt/quantum_trader"
EnvironmentFile=/etc/quantum/ai-engine.env
ExecStart=/opt/quantum/venvs/ai-engine/bin/uvicorn microservices.ai_engine.service:app --host 0.0.0.0 --port 8001 --log-level info

# Restart policy
Restart=always
RestartSec=2

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# Resource limits
MemoryMax=2G
MemoryHigh=1.8G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-ai-engine

[Install]
WantedBy=quantum-ai.target
