# ops/systemd/native/quantum-market-proof.service
[Unit]
Description=Quantum Market Data Proof (E2E streams test)
After=redis-server.service quantum-market-publisher.service
Wants=redis-server.service quantum-market-publisher.service

[Service]
Type=oneshot
RemainAfterExit=no
WorkingDirectory=/home/qt/quantum_trader
ExecStartPre=/bin/systemctl is-active --quiet quantum-market-publisher.service
ExecStart=/bin/bash /home/qt/quantum_trader/ops/tests/test_market_data_e2e.sh
TimeoutStartSec=120
User=root

[Install]
WantedBy=multi-user.target
