# ops/systemd/native/quantum-execution.service
[Unit]
Description=Quantum Trader - Execution (native with governor protection)
After=network.target redis-server.service quantum-ai-engine.service
Requires=redis-server.service
PartOf=quantum-exec.target

[Service]
Type=simple
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
EnvironmentFile=-/etc/quantum/execution.env
Environment="PYTHONPATH=/home/qt/quantum_trader"

# Using main_v2.py with Governor Gate
ExecStart=/opt/quantum/venvs/execution/bin/python -m uvicorn microservices.execution.main_v2:app --host 127.0.0.1 --port 8003

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-execution

[Install]
WantedBy=quantum-exec.target
