# ops/systemd/native/quantum-ai-engine-proof.service
[Unit]
Description=Quantum AI Engine Publish Proof (Real signals test)
After=redis-server.service quantum-ai-engine.service
Wants=redis-server.service quantum-ai-engine.service

[Service]
Type=oneshot
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/qt/quantum_trader

# Environment variables
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=AI_PROOF_MODE=activity
Environment=AI_PROOF_MAX_AGE_SECONDS=600
Environment=AI_PROOF_WAIT_SECONDS=60

# Pre-check and execution
ExecStartPre=/bin/systemctl is-active --quiet quantum-ai-engine.service
ExecStart=/bin/bash /home/qt/quantum_trader/ops/tests/test_ai_engine_publish.sh
TimeoutStartSec=90
User=root

[Install]
WantedBy=multi-user.target
