# deploy/quantum-offline-evaluator.service
[Unit]
Description=Quantum Trader - Exit Brain v1 Offline Evaluator (Phase 5)
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=oneshot
User=qt
Group=qt
WorkingDirectory=/home/qt/quantum_trader
Environment="PYTHONPATH=/home/qt/quantum_trader"
Environment="REDIS_HOST=127.0.0.1"
Environment="REDIS_PORT=6379"
Environment="REDIS_DB=0"
Environment="EVAL_WINDOW_HOURS=4"
Environment="EVAL_HORIZON_SECONDS=14400"
ExecStart=/home/qt/quantum_trader_venv/bin/python3 -m microservices.exit_brain_v1.run_offline_evaluation
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-offline-eval

[Install]
WantedBy=multi-user.target
