# deploy/quantum-offline-evaluator.timer
[Unit]
Description=Run Exit Brain v1 Offline Evaluator every 4 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target
