# systemd/units/quantum-rl-feedback-v2.service
[Unit]
Description=Quantum Trader - RL Feedback V2 (AI Client)
After=network.target quantum-rl-sizer.service

[Service]
Type=simple
WorkingDirectory=/mnt/c/quantum_trader

ExecStart=/usr/bin/python3 rl_feedback_v2_daemon.py

MemoryHigh=230M
MemoryMax=256M
CPUQuota=40%
TasksMax=64

Restart=always
RestartSec=5s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-rl-feedback-v2

# Security constraints removed for local dev environment

[Install]
WantedBy=quantum-trader.target
