# systemd/units/quantum-nginx-proxy.service
[Unit]
Description=Quantum Trader - NGINX Reverse Proxy (Infrastructure)
After=network.target
PartOf=quantum-trader.target

[Service]
Type=forking
User=www-data
Group=www-data

PIDFile=/run/nginx.pid
ExecStartPre=/usr/sbin/nginx -t -c /etc/quantum/nginx.conf
ExecStart=/usr/sbin/nginx -c /etc/quantum/nginx.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID

MemoryHigh=230M
MemoryMax=256M
CPUQuota=50%

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-nginx

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/nginx /var/lib/nginx /run

[Install]
WantedBy=quantum-trader.target
