# systemd/units/quantum-trade-intent-consumer.service
[Unit]
Description=Quantum Trader - Trade Intent Consumer (AI Client)
After=network.target quantum-redis.service
Requires=quantum-redis.service
PartOf=quantum-trader.target

[Service]
Type=simple
User=quantum-trade-intent-consumer
Group=quantum-trade-intent-consumer
WorkingDirectory=/opt/quantum

EnvironmentFile=/etc/quantum/trade-intent-consumer.env

ExecStart=/opt/quantum/venvs/trade-intent-consumer/bin/python \
    trade_intent_runner.py

MemoryHigh=480M
MemoryMax=512M
CPUQuota=50%
TasksMax=128

Restart=always
RestartSec=10s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-trade-intent-consumer

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/quantum /var/log/quantum
ReadOnlyPaths=/opt/quantum

[Install]
WantedBy=quantum-trader.target
