# systemd/units/quantum-rl-sizer.service
[Unit]
Description=Quantum Trader - RL Position Sizing Agent (Model Server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/c/quantum_trader

ExecStart=/usr/bin/python3 -m microservices.rl_sizing_agent.pnl_feedback_listener

# Resource limits (MODEL SERVER - MEDIUM)
MemoryHigh=350M
MemoryMax=384M
CPUQuota=100%
TasksMax=128

# Restart policy
Restart=always
RestartSec=5s
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-rl-sizer

# Logging
# (Security constraints removed for local dev environment)

[Install]
WantedBy=quantum-trader.target
