# systemd/units/quantum-frontend.service
[Unit]
Description=Quantum Trader - Main Frontend (Infrastructure)
After=network.target quantum-redis.service
PartOf=quantum-trader.target

[Service]
Type=simple
User=quantum-frontend
Group=quantum-frontend
WorkingDirectory=/opt/quantum/frontend

EnvironmentFile=/etc/quantum/frontend.env

ExecStart=/usr/bin/node server.js

MemoryHigh=230M
MemoryMax=256M
CPUQuota=50%
TasksMax=64

Restart=always
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-frontend

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/quantum

[Install]
WantedBy=quantum-trader.target
