# systemd/units/quantum-strategy-ops.service
[Unit]
Description=Quantum Trader - Strategy Operations (Model Server)
After=network.target quantum-redis.service
Requires=quantum-redis.service
PartOf=quantum-trader.target

[Service]
Type=simple
User=quantum-strategy-ops
Group=quantum-strategy-ops
WorkingDirectory=/opt/quantum/microservices/strategy_operations

EnvironmentFile=/etc/quantum/strategy-ops.env

ExecStart=/opt/quantum/venvs/strategy-ops/bin/python strategy_ops.py

# Resource limits (MODEL SERVER - MEDIUM)
MemoryHigh=350M
MemoryMax=384M
CPUQuota=100%
TasksMax=128

# Restart policy
Restart=always
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-strategy-ops

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/quantum /var/log/quantum
ReadOnlyPaths=/opt/quantum

[Install]
WantedBy=quantum-trader.target
