# systemd/units/quantum-risk-safety.service
[Unit]
Description=Quantum Trader - Risk Safety Monitor (AI Client)
After=network.target quantum-redis.service quantum-risk-brain.service
Requires=quantum-redis.service
PartOf=quantum-trader.target

[Service]
Type=simple
User=quantum-risk-safety
Group=quantum-risk-safety
WorkingDirectory=/opt/quantum/microservices/risk_safety

EnvironmentFile=/etc/quantum/risk-safety.env

ExecStart=/opt/quantum/venvs/risk-safety/bin/python main.py

MemoryHigh=230M
MemoryMax=256M
CPUQuota=30%
TasksMax=64

Restart=always
RestartSec=5s
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantum-risk-safety

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/quantum /var/log/quantum
ReadOnlyPaths=/opt/quantum

[Install]
WantedBy=quantum-trader.target
